# docker/cron/crontab/docker
15 0 * * * /usr/local/bin/php /var/www/artisan import >> /dev/null
# * * * * * /bin/echo "foobar" >> /var/log/cron.log
